# deploy/temple-solver.service
[Unit]
Description=POE2 Temple Solver API
After=network.target

[Service]
Type=exec
User=www-data
Group=www-data
WorkingDirectory=/opt/temple-solver/solver-python
Environment="PATH=/opt/temple-solver/venv/bin"
Environment="MAX_CONCURRENT_SOLVES=2"
Environment="MAX_QUEUE_SIZE=10"
Environment="RATE_LIMIT_SECONDS=30"
Environment="MAX_SOLVE_TIME=60"
ExecStart=/opt/temple-solver/venv/bin/gunicorn -c /opt/temple-solver/deploy/gunicorn.conf.py server_prod:app
ExecReload=/bin/kill -s HUP $MAINPID
Restart=always
RestartSec=5

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/log/temple-solver
PrivateTmp=true

[Install]
WantedBy=multi-user.target
